# docker/elogen/crontab
0 19 * * * cncnet /app/elogen --output-dir /data -l info